# apps/admin/ops/liu-blog-admin-web.service
[Unit]
Description=Xinghui Blog server administration frontend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=blog
Group=blog
WorkingDirectory=/srv/xinghui-blog-admin/source/.next/standalone
Environment=NODE_ENV=production
Environment=NEXT_TELEMETRY_DISABLED=1
Environment=HOSTNAME=127.0.0.1
Environment=PORT=3001
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5s
TimeoutStopSec=20s
MemoryMax=512M
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/xinghui-blog-admin /var/lib/xinghui-blog-admin

[Install]
WantedBy=multi-user.target
